# deploy/supervisord-uwsgi.conf
[program:uwsgi]
command=uwsgi --die-on-term --ini /code/gaa/deploy/uwsgi.ini
directory=/code/gaa
autostart=true
autorestart=true
redirect_stderr=true
